# scripts/systemd/dashboard-pull.service
[Unit]
Description=Dashboard Pull (git pull + build + restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clawd-dashboard-update.sh /etc/clawd/template.env

[Install]
WantedBy=multi-user.target
