# scripts/openclaw-dashboard-watchdog.service
[Unit]
Description=OpenClaw Dashboard availability watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=clawdbot
Group=clawdbot
WorkingDirectory=/home/clawdbot/apps/dashboard
Environment=OPENCLAW_HOME=/home/clawdbot/.openclaw
Environment=DASHBOARD_ENV_PATH=/etc/clawd/template.env
ExecStart=/usr/bin/env bash /home/clawdbot/apps/dashboard/scripts/openclaw-dashboard-watchdog.sh
