# scripts/openclaw-dashboard-watchdog.timer
[Unit]
Description=Run dashboard watchdog every minute

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=openclaw-dashboard-watchdog.service

[Install]
WantedBy=timers.target
