# scripts/openclaw-task-orchestrator.service
[Unit]
Description=OpenClaw Task Orchestrator (ticket auto-pickup worker)
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=clawdbot
Group=clawdbot
WorkingDirectory=/home/clawdbot/apps/dashboard
Environment=HOME=/home/clawdbot
Environment=DASHBOARD_ENV_PATH=/etc/clawd/template.env
ExecStart=/usr/bin/env node /home/clawdbot/apps/dashboard/scripts/task-orchestrator.js
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
